# src/main/assembly/dist/install/dd-vault-ingest-flow.service
[Unit]
Description=Vault Ingest Flow Service

[Service]
ExecStart=/opt/dans.knaw.nl/dd-vault-ingest/bin/dd-vault-ingest server /etc/opt/dans.knaw.nl/dd-vault-ingest/config.yml
# Java returns 143 even if the SIGTERM was handled correctly.
SuccessExitStatus=143
Environment="JAVA_OPTS=-Xms64m"
Environment="LC_ALL=en_US.UTF-8"

User=dd-vault-ingest
Group=dd-vault-ingest

[Install]
WantedBy=multi-user.target
